# bind/named.conf
[program:named]
command=bash -c '/usr/sbin/named -f -c /etc/bind/named.conf -d 5'
startsecs=0
